# deployment/lokal.service
[Unit]
Description=LoKal Educational AI Backend
After=network.target ollama.service

[Service]
Type=simple
User=lokal
Group=lokal
WorkingDirectory=/home/lokal/Desktop/loKal
Environment="PATH=/home/lokal/Desktop/loKal/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/lokal/Desktop/loKal/venv/bin/gunicorn --config /home/lokal/Desktop/loKal/deployment/gunicorn.conf.py lokal_backend.wsgi:application
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/lokal/Desktop/loKal/logs/gunicorn.log
StandardError=append:/home/lokal/Desktop/loKal/logs/gunicorn_error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
